# rootfs/etc/s6-overlay/s6-rc.d/svc-krusader-ready/finish
#!/command/with-contenv bash
# -----------------------------------------------------------------------------
# svc-krusader-ready/finish — runs when the service is brought down, i.e. on
# container stop. Quit krusader CLEANLY (kquitapp) so it writes its "save on
# exit" UI state (krusaderrc, with [Startup] User Save Settings=true) BEFORE s6
# SIGKILLs the desktop. This is the only chance to persist a status bar the user
# hid via the View menu — `docker stop` otherwise kills krusader before it saves.
# Best-effort + capped; experimental (Krusader 2.8.1 may still ignore the key).
# -----------------------------------------------------------------------------
_addr="$(cat /tmp/krusader-dbus-address 2>/dev/null)"
[ -n "${_addr}" ] && export DBUS_SESSION_BUS_ADDRESS="${_addr}"

kquitapp6 krusader 2>/dev/null \
    || kquitapp5 krusader 2>/dev/null \
    || qdbus6 org.kde.krusader /MainApplication quit 2>/dev/null \
    || qdbus   org.kde.krusader /MainApplication quit 2>/dev/null \
    || true

# Give krusader a moment to write krusaderrc before the s6 SIGKILL sweep.
sleep 2
